# systemd/templates/poc.service.j2
[Unit]
Description=POC Daemon

[Service]
ExecStart=/bin/bash /poc/poc.sh
User=root

[Install]
WantedBy=multi-user.target